# ansible/files/supabase_admin_agent_config/supabase-admin-agent_salt.service
[Unit]
Description=Configuration management via supabase-admin-agent salt
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/supabase-admin-agent/supabase-admin-agent --config /opt/supabase-admin-agent/config.yaml salt --apply --store-result
User=supabase-admin-agent
Group=supabase-admin-agent
StandardOutput=journal
StandardError=journal
StateDirectory=supabase-admin-agent
CacheDirectory=supabase-admin-agent

# Security hardening
PrivateTmp=true

[Install]
WantedBy=multi-user.target
